# packaging/photo-frame-display.service
[Unit]
Description=Photo Frame DRM Display
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/photo-frame-display
EnvironmentFile=-/etc/photo-frame/display.env
Restart=on-failure
User=photo-frame
Group=video
RuntimeDirectory=photo-frame

[Install]
WantedBy=multi-user.target
